# appliance-root/opt/napp/bin/noitd-ctlr
#!/usr/bin/python
#
# chkconfig: 345 25 75
# description: A daemon that starts and stops noitd by file presense


import sys, os, fcntl, time, struct

def set_initial_state():
    if not os.system("pgrep '^noitd$'"):
        f = open('/opt/napp/etc/noit.run', 'w')
        f.close()
        return True
    else:
        return False
    

def main(cmd):

    if cmd == 'start':
        print('Starting...')
    elif cmd == 'stop':
        if not os.path.isfile('/var/run/noitd-ctrl.lock'):
            print('not running...')
            os._exit(-1)
        try:
            lock = open('/var/run/noitd-ctrl.lock', 'r')
            pid = lock.read()
            lock.close()
            os.kill(int(pid), 3)
            os.unlink('/var/run/noitd-ctrl.lock')
            os._exit(0)
        except IOError:
            os._exit(-1)
    else:
        print('%s not understood' % (cmd))
        os._exit(-1)

    pid = os.fork()
    if (pid != 0):
        os._exit(0)
    os.setsid()
    pid = os.fork()
    if (pid != 0):
        os._exit(0)

    lock = os.open('/var/run/noitd-ctrl.lock', os.O_RDWR | os.O_CREAT)
    lockdata = struct.pack('hhllhh', fcntl.F_WRLCK, 0, 0, 0, 0, 0)
    try:
        rv = fcntl.fcntl(lock, fcntl.F_SETLK, lockdata)
    except IOError:
        print("Already running...")
        os._exit(-1)
    os.write(lock, "%d" % (os.getpid()))

    noit_on = set_initial_state()

    while True:
        if os.path.isfile('/opt/napp/etc/noit.run'):
            if not noit_on:
                print("starting noitd")
                os.system("/etc/rc.d/init.d/noitd start")
                noit_on = True
        else:
            if noit_on:
                print("stopping noitd")
                os.system("/etc/rc.d/init.d/noitd stop")
                noit_on = False
        time.sleep(2)

if __name__ == "__main__":
    main((lambda x:sys.argv[1], lambda x:'command')[len(sys.argv) != 2](0))
