# appliance-root/etc/cron.d/napp
* * * * * root /opt/napp/bin/issue-refresh > /dev/null 2>&1
0 * * * * root /opt/napp/bin/crl-refresh > /dev/null 2>&1
* * * * * root /opt/napp/bin/perform-updates > /dev/null 2>&1
